# setup/crontab
@reboot sleep 30 && /home/pi/calcifair/setup/start.sh &
0 2 * * * pkill calcifair-main
0 7 * * * if [ -z `pgrep calcifair` ]; then /usr/bin/python3 /home/pi/calcifair/calcifair.py; fi